# systemd/yeoni-farm.service
[Unit]
Description=Yeoni Farm static website
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/farmer/data/yeoni
ExecStart=/home/farmer/.local/bin/caddy run --config /home/farmer/data/yeoni/Caddyfile
ExecReload=/home/farmer/.local/bin/caddy reload --config /home/farmer/data/yeoni/Caddyfile
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
